# supervisor/unsafepay.conf.sample
[program:unsafepay]
command=/path/venv/bin/python /path/unsafepay/bot.py
director=/path/unsafepay/
environment=LNDDIR=/home/lightning/.lnd
user=lightning
